# dse/datastax-agent/bin/com.datastax.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.datastax.agent</string>
	<key>RunAtLoad</key>
	<true/>
        <key>KeepAlive</key>
        <false/>
        <key>AbandonProcessGroup</key>
        <false/>
        <key>WorkingDirectory</key>
        <string>/Users/crupib/dse/datastax-agent</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/crupib/dse/datastax-agent/bin/datastax-agent</string>
                <string>-f</string>
	</array>
        <key>StandardErrorPath</key>
        <string>/Users/crupib/dse/logs/datastax-agent/agent_error.log</string>
        <key>StandardOutPath</key>
        <string>/Users/crupib/dse/logs/datastax-agent/agent_output.log</string>
</dict>
</plist>
